# api/tilestache.upstart
# Starting conditions
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

chdir /home/ubuntu/
setuid nobody

script
exec /home/ubuntu/venv/bin/uwsgi --master --socket=127.0.0.1:8000 --eval 'import TileStache; application = TileStache.WSGITileServer("/home/ubuntu/tilestache.cfg")' --enable-threads -p 4
end script



